# supervisor/supervisord.conf
{% for appname,params in pillar['supervisor']['apps'].iteritems() -%}
[program:{{ appname }}]
command={{ params['command'] }}
{% endfor %}
