# xvnc/xvnc.conf
[program:xvnc]
command=/usr/bin/x11vnc -reopen -forever
autostart=true
autorestart=true
